# crontab
* * * * * /.venv/bin/python /app/main.py
