# crontab
*/20 * * * * * * /.venv/bin/python /app/main.py
